# conf/retroarch.cron
#This script will update the core index every 5 minutes so that they are available in retroarch
*/5 * * * * root __FINAL_PATH__/indexer.sh > /dev/null 2>&1
