# conf/retroarch.cron
*/5 * * * * root __FINAL_PATH__/indexer.sh > /dev/null 2>&1